# systemd-user/prez-pkglog.service
[Unit]
Description=Prez Package Logger (Downloads Monitor)
Documentation=man:prez-pkglog(1)
After=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/prez-pkglog daemon --scope user
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
